# deploy/pip.service
# Pip - systemd Service File
#
# Installation:
#   sudo cp deploy/pip.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable pip
#   sudo systemctl start pip
#
# Logs:
#   sudo journalctl -u pip -f

[Unit]
Description=Pip - AI Accounting Assistant
Documentation=https://github.com/IAMSamuelRodda/pip
After=network.target

[Service]
Type=simple
User=pip
Group=pip

# Working directory (adjust path as needed)
WorkingDirectory=/opt/pip

# Node.js executable (adjust if using nvm or different path)
ExecStart=/usr/bin/node packages/server/dist/index.js

# Environment file with secrets
EnvironmentFile=/opt/pip/.env

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3000

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/pip/data

# Resource limits
MemoryMax=512M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pip

[Install]
WantedBy=multi-user.target
